# deploy/systemd/esptool-sync.timer
[Unit]
Description=Run EastMesh firmware sync every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
Unit=esptool-sync.service

[Install]
WantedBy=timers.target
